# com.local.proxy-router.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.local.proxy-router</string>

  <key>ProgramArguments</key>
  <array>
    <!-- Update this path after building -->
    <string>/usr/local/bin/proxy-router</string>
    <string>-config</string>
    <string>/usr/local/etc/proxy-router/config.json</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/proxy-router.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/proxy-router.err</string>
</dict>
</plist>
